# conf/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# make sure fancontrol starts by running this script
fancontroller

# disable bluetooth (possible to reenable later)
rfkill block bluetooth

# set charge limits for battery 0
batt-charge-limit 60 80 0

exit 0
