# conf/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# rewritefs
#rewritefs -o config=/mnt/home/jure/.config/rewritefs,allow_other /mnt/home/jure /home/jure

# make sure fancontrol starts by running this script
fancontroller babysitter 2>>/var/log/fancontroller.log &

# disable bluetooth (possible to reenable later)
rfkill block bluetooth

# set charge limits for battery 0
batt-charge-limit 60 80 0

# move some junk to RAM
rammer > /dev/null 2>>/var/log/rammer.log

exit 0
